# proxyProject/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/local/bin/gunicorn --bind unix:/application/gunicorn.sock proxyProject.wsgi:application
directory=/application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log